# conf/supervisord-mysqld.conf
[program:mysqld]
command=/scripts/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
